# scripts/mdm/macos/org.hol.guard.machine-health.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>org.hol.guard.machine-health</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Library/Application Support/HOL Guard/hol-guard/hol-guard</string>
    <string>mdm</string>
    <string>health-report</string>
    <string>--scope</string>
    <string>machine</string>
    <string>--json</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key><string>/var/root</string>
    <key>PATH</key><string>/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>TMPDIR</key><string>/var/tmp</string>
  </dict>
  <key>RunAtLoad</key><true/>
  <key>StartInterval</key><integer>300</integer>
  <key>ProcessType</key><string>Background</string>
  <key>StandardOutPath</key><string>/Library/Logs/HOL Guard/machine-health.log</string>
  <key>StandardErrorPath</key><string>/Library/Logs/HOL Guard/machine-health-error.log</string>
</dict>
</plist>
